# piday2025posters/display-hat-gallery.service
[Unit]
Description=Display HAT Mini Image Gallery
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/piday2025/piday2025posters
ExecStart=/home/pi/piday2025/piday2025posters/startup-script.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
EOF